# django.ini
[program:django]
command=/Users/achilles_xushy/.local/share/virtualenvs/django/bin/gunicorn -c gunicorn_config.py achilles_book.wsgi:application
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
numprocs_start=1
stopwaitsecs=5
user=achilles_xushy
directory=/Users/achilles_xushy/PycharmProjects/achilles_book
stdout_logfile=/Users/achilles_xushy/PycharmProjects/config_and_log/supervisor/gunicorn_django_access.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/Users/achilles_xushy/PycharmProjects/config_and_log/supervisor/gunicorn_django_error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
