# gunicorn.ini
[program:gunicorn]
command=/Users/achilles_xushy/.pyenv/versions/anaconda3-4.0.0/bin/gunicorn -c gunicorn_config.py achilles_book.wsgi:application
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
numprocs_start=1
stopwaitsecs=12
user=achilles_xushy
directory=/Users/achilles_xushy/Desktop/python2/achilles_book
stdout_logfile=/Users/achilles_xushy/Desktop/python2/7po_server/supervisor/gunicorn.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/Users/achilles_xushy/Desktop/python2/7po_server/supervisor/gunicornerr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
